# systemd/cert-checker.service.in
[Unit]
Description=Start the Cert-Checker
Requires=dbus.service

[Service]
User=@SERVICE_USER@
Group=@SERVICE_GROUP@
SmackProcessLabel=System
ExecStart=@BIN_DIR@/cert-checker
Sockets=@SERVICE_NAME@.socket
