# systemd/cert-checker-popup.socket.in
[Unit]
Description=Cert-checker popup socket

[Socket]
ListenStream=@POPUP_STREAM@
SocketMode=0777

Service=@SERVICE_NAME@-popup.service

[Install]
WantedBy=sockets.target
